# sources/api/debian/ostdlabs-etracker-customer-api.cron.d
0 3 * * * ${customer} usr/share/ostdlabs/etracker/${customer}/${customer}-api/sql/control/backup.sh

